# root/etc/s6-overlay/s6-rc.d/init-mod-trusted-certs-postinstall/up
/etc/s6-overlay/s6-rc.d/init-mod-trusted-certs-postinstall/run