# root/etc/s6-overlay/s6-rc.d/init-mod-trusted-certs-install/run
#!/command/with-contenv bash

echo "ca-certificates" >> /mod-repo-packages-to-install.list