# root/etc/s6-overlay/s6-rc.d/init-mod-trusted-certs-install/up
/etc/s6-overlay/s6-rc.d/init-mod-trusted-certs-install/run